# autostart/wallpaper.desktop
[Desktop Entry]
Exec=/home/vagrant/.config/autostart/set_wallpaper.sh
Name=Set wallpaper
Type=Application
Version=1.0